# trunk/src/plugins/nagios2/etc/cron.d/openqrm-nagios2-automap.cron
*/5 * * * * root @@OPENQRM_SERVER_BASE_DIR@@/openqrm/plugins/nagios2/bin/openqrm-nagios-manager map
