# trunk/src/plugins/puppet/etc/cron.d/openqrm-puppetm-config.cron
* * * * * root @@OPENQRM_SERVER_BASE_DIR@@/openqrm/plugins/puppet/bin/openqrm-puppetm-config
